# frontend/deploy/systemd/dpr-web.service
[Unit]
Description=DPR Web Next.js Service
After=network.target

[Service]
Type=simple
User=teksys
Group=teksys
WorkingDirectory=/opt/dpr-teksys/frontend/dpr-web
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3000
ExecStart=/opt/dpr-teksys/frontend/dpr-web/scripts/run-prod.sh
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
